# infra/systemd/chat-app-idle.service
[Unit]
Description=Chat-App idle check — scale the ASG to 0 after inactivity
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=ACTIVITY_FILE=/opt/chat-app/infra/compose/state/last-activity
Environment=IDLE_MINUTES=30
Environment=AWS_REGION=ap-south-1
ExecStart=/opt/chat-app/infra/prod/idle-check.sh
